# rpi/athocs-sensor.cron

# athocs sensor
*/5 * * * * . ~/.profile && cd "$ATHOCS_BASE_DIR/athocs-sensor/rpi" && /usr/bin/python3 "athocs-sensor.py" >"$ATHOCS_BASE_DIR/log/sensor-errors.log" 2>&1

